# deploy/server/energy-server.service
[Unit]
Description=Synctacles Energy Price Server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=synctacles
Group=synctacles
WorkingDirectory=/opt/energy-server
ExecStart=/opt/energy-server/energy-server
Restart=always
RestartSec=10
Environment=DATABASE_URL=postgresql://energy:CHANGEME@localhost:5432/energy_prod?sslmode=disable
Environment=LISTEN_ADDR=:8080
Environment=DEBUG_MODE=false
StandardOutput=journal
StandardError=journal
SyslogIdentifier=energy-server

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/energy-server

[Install]
WantedBy=multi-user.target
