# systemd/synctacles-tennet.timer.template
[Unit]
Description={{BRAND_NAME}} TenneT Timer (5 min, acceptance API)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
