# systemd/synctacles-tennet.service
[Unit]
Description=SYNCTACLES TenneT Collector (Rate Limited)
After=network.target

[Service]
Type=oneshot
User=synctacles
WorkingDirectory=/opt/synctacles/app
ExecStart=/opt/synctacles/venv/bin/python3 /opt/synctacles/app/sparkcrawler_db/collectors/sparkcrawler_tennet_ingestor.py
StandardOutput=journal
StandardError=journal
