# systemd/synctacles-collector.service.template
[Unit]
Description={{BRAND_NAME}} Data Collector
After=network.target postgresql.service

[Service]
Type=oneshot
User={{SERVICE_USER}}
WorkingDirectory={{INSTALL_PATH}}/app
ExecStart={{INSTALL_PATH}}/app/scripts/run_collectors.sh
StandardOutput=journal
StandardError=journal
