# systemd/synctacles-tennet.timer
[Unit]
Description=SYNCTACLES TenneT Timer (5 min, acceptance API)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
